# debianrfs/var/lib/dpkg/info/bb-customizations.list
/.
/etc
/etc/avahi
/etc/avahi/services
/etc/avahi/services/ssh.service
/etc/bbb.io
/etc/bbb.io/examples
/etc/bbb.io/examples/mplayer
/etc/bbb.io/templates
/etc/bbb.io/templates/dev
/etc/bbb.io/templates/dev/bone
/etc/bbb.io/templates/dev/bone/counter
/etc/bbb.io/templates/dev/bone/counter/README.md
/etc/bbb.io/templates/dev/bone/i2c
/etc/bbb.io/templates/dev/bone/i2c/README.md
/etc/bbb.io/templates/dev/bone/pwm
/etc/bbb.io/templates/dev/bone/pwm/README.md
/etc/bbb.io/templates/dev/bone/spi
/etc/bbb.io/templates/dev/bone/spi/README.md
/etc/bbb.io/templates/dev/bone/uart
/etc/bbb.io/templates/dev/bone/uart/README.md
/etc/bbb.io/templates/nginx
/etc/bbb.io/templates/nginx/BeaglePlay-io.html
/etc/bbb.io/templates/nginx/Node-RED-ui.html
/etc/bbb.io/templates/nginx/Node-RED.html
/etc/bbb.io/templates/nginx/VSCode.html
/etc/bbb.io/templates/nginx/nginx-autoindex
/etc/bbb.io/templates/xfce4
/etc/bbb.io/templates/xfce4/xfce4-desktop.xml
/etc/kernel
/etc/kernel/postinst.d
/etc/kernel/postinst.d/zz-uenv_txt
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/10-of-symlink.rules
/etc/udev/rules.d/80-gpio-noroot.rules
/etc/udev/rules.d/80-i2c-noroot.rules
/etc/udev/rules.d/81-pwm-noroot.rules
/etc/udev/rules.d/82-gpio-config-pin.rules
/etc/udev/rules.d/83-counter-noroot.rules
/etc/udev/rules.d/85-gpio-noroot.rules
/etc/udev/rules.d/86-remoteproc-noroot.rules
/etc/udev/rules.d/86-rpmsg-noroot.rules
/etc/udev/rules.d/87-iio-noroot.rules
/etc/udev/rules.d/88-leds-noroot.rules
/etc/udev/rules.d/uio-pruss.rules
/etc/udev/rules.d/uio.rules
/lib
/lib/firmware
/lib/firmware/am335x-pru0-fw.sleep
/lib/firmware/am335x-pru1-fw.sleep
/lib/firmware/am57xx-pru1_0-fw.sleep
/lib/firmware/am57xx-pru1_1-fw.sleep
/lib/firmware/am57xx-pru2_0-fw.sleep
/lib/firmware/am57xx-pru2_1-fw.sleep
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bb-symlinks.service
/usr
/usr/sbin
/usr/sbin/bb-symlinks
/usr/share
/usr/share/backgrounds
/usr/share/backgrounds/bbb.io
/usr/share/backgrounds/bbb.io/beaglebg.jpg
/usr/share/backgrounds/bbb.io/beagleboard-logo.svg
/usr/share/doc
/usr/share/doc/bb-customizations
/usr/share/doc/bb-customizations/changelog.Debian.gz
/usr/share/doc/bb-customizations/copyright
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/ti_pru_firmware
