# debianrfs/var/lib/dpkg/info/bb-usb-gadgets.list
/.
/etc
/etc/bbb.io
/etc/bbb.io/templates
/etc/bbb.io/templates/bb-boot
/etc/bbb.io/templates/beaglebg.jpg
/etc/bbb.io/templates/beagleboard-logo.svg
/etc/bbb.io/templates/eth0-DHCP.network
/etc/bbb.io/templates/eth1-DHCP.network
/etc/bbb.io/templates/systemd-networkd-wait-online.override.conf
/etc/bbb.io/templates/usb0-DHCP.network
/etc/bbb.io/templates/usb0-DHCPServer.network
/etc/bbb.io/templates/usb1-DHCP.network
/etc/bbb.io/templates/usb1-DHCPServer.network
/etc/gt
/etc/gt/templates
/etc/gt/templates/bbb-acm-mass_storage-ncm-rndis.scheme
/etc/gt/templates/bbb-acm-ncm-rndis.scheme
/etc/gt/templates/bbb-acm-rndis.scheme
/etc/gt/templates/bbb-acm.scheme
/etc/gt/templates/bbb-ncm-rndis.scheme
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bb-usb-gadgets.service
/usr
/usr/bin
/usr/bin/autoconfigure_usb0.sh
/usr/bin/autoconfigure_usb1.sh
/usr/bin/bb-start-acm-ncm-rndis-old-gadget
/usr/bin/bb-start-usb-gadgets
/usr/bin/bb_dnsmasq_config.sh
/usr/bin/bb_generate_mac.sh
/usr/share
/usr/share/doc
/usr/share/doc/bb-usb-gadgets
/usr/share/doc/bb-usb-gadgets/changelog.Debian.gz
/usr/share/doc/bb-usb-gadgets/copyright
