# debianrfs/var/lib/dpkg/info/cockpit-ws.list
/.
/etc
/etc/cockpit
/etc/cockpit/ws-certs.d
/etc/pam.d
/etc/pam.d/cockpit
/lib
/lib/arm-linux-gnueabihf
/lib/arm-linux-gnueabihf/security
/lib/arm-linux-gnueabihf/security/pam_cockpit_cert.so
/lib/arm-linux-gnueabihf/security/pam_ssh_add.so
/lib/systemd
/lib/systemd/system
/lib/systemd/system/cockpit-motd.service
/lib/systemd/system/cockpit-wsinstance-http.service
/lib/systemd/system/cockpit-wsinstance-http.socket
/lib/systemd/system/cockpit-wsinstance-https-factory.socket
/lib/systemd/system/cockpit-wsinstance-https-factory@.service
/lib/systemd/system/cockpit-wsinstance-https@.service
/lib/systemd/system/cockpit-wsinstance-https@.socket
/lib/systemd/system/cockpit.service
/lib/systemd/system/cockpit.socket
/lib/systemd/system/system-cockpithttps.slice
/usr
/usr/lib
/usr/lib/cockpit
/usr/lib/cockpit/cockpit-certificate-ensure
/usr/lib/cockpit/cockpit-certificate-helper
/usr/lib/cockpit/cockpit-client
/usr/lib/cockpit/cockpit-client.ui
/usr/lib/cockpit/cockpit-desktop
/usr/lib/cockpit/cockpit-session
/usr/lib/cockpit/cockpit-tls
/usr/lib/cockpit/cockpit-ws
/usr/lib/cockpit/cockpit-wsinstance-factory
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/cockpit-tempfiles.conf
/usr/share
/usr/share/cockpit
/usr/share/cockpit/branding
/usr/share/cockpit/branding/arch
/usr/share/cockpit/branding/arch/branding.css
/usr/share/cockpit/branding/debian
/usr/share/cockpit/branding/debian/branding.css
/usr/share/cockpit/branding/default
/usr/share/cockpit/branding/default/apple-touch-icon.png
/usr/share/cockpit/branding/default/bg-plain.jpg
/usr/share/cockpit/branding/default/brand-large.png
/usr/share/cockpit/branding/default/branding.css
/usr/share/cockpit/branding/default/favicon.ico
/usr/share/cockpit/branding/default/logo.png
/usr/share/cockpit/branding/kubernetes
/usr/share/cockpit/branding/kubernetes/branding.css
/usr/share/cockpit/branding/registry
/usr/share/cockpit/branding/registry/branding.css
/usr/share/cockpit/motd
/usr/share/cockpit/motd/inactive.motd
/usr/share/cockpit/motd/update-motd
/usr/share/cockpit/static
/usr/share/cockpit/static/fonts
/usr/share/cockpit/static/fonts/RedHatDisplay-Black.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-BlackItalic.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-Bold.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-BoldItalic.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-Italic.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-Medium.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-MediumItalic.woff2
/usr/share/cockpit/static/fonts/RedHatDisplay-Regular.woff2
/usr/share/cockpit/static/fonts/RedHatMono-Bold.woff2
/usr/share/cockpit/static/fonts/RedHatMono-BoldItalic.woff2
/usr/share/cockpit/static/fonts/RedHatMono-Italic.woff2
/usr/share/cockpit/static/fonts/RedHatMono-Medium.woff2
/usr/share/cockpit/static/fonts/RedHatMono-MediumItalic.woff2
/usr/share/cockpit/static/fonts/RedHatMono-Regular.woff2
/usr/share/cockpit/static/fonts/RedHatText-Bold.woff2
/usr/share/cockpit/static/fonts/RedHatText-BoldItalic.woff2
/usr/share/cockpit/static/fonts/RedHatText-Italic.woff2
/usr/share/cockpit/static/fonts/RedHatText-Medium.woff2
/usr/share/cockpit/static/fonts/RedHatText-MediumItalic.woff2
/usr/share/cockpit/static/fonts/RedHatText-Regular.woff2
/usr/share/cockpit/static/login.css
/usr/share/cockpit/static/login.html
/usr/share/cockpit/static/login.js
/usr/share/cockpit/static/manifest.json
/usr/share/cockpit/static/po.cs.js
/usr/share/cockpit/static/po.de.js
/usr/share/cockpit/static/po.es.js
/usr/share/cockpit/static/po.fi.js
/usr/share/cockpit/static/po.fr.js
/usr/share/cockpit/static/po.he.js
/usr/share/cockpit/static/po.it.js
/usr/share/cockpit/static/po.ja.js
/usr/share/cockpit/static/po.ka.js
/usr/share/cockpit/static/po.ko.js
/usr/share/cockpit/static/po.nb_NO.js
/usr/share/cockpit/static/po.nl.js
/usr/share/cockpit/static/po.pl.js
/usr/share/cockpit/static/po.pt_BR.js
/usr/share/cockpit/static/po.ru.js
/usr/share/cockpit/static/po.sk.js
/usr/share/cockpit/static/po.sv.js
/usr/share/cockpit/static/po.tr.js
/usr/share/cockpit/static/po.uk.js
/usr/share/cockpit/static/po.zh_CN.js
/usr/share/cockpit/static/po.zh_TW.js
/usr/share/doc
/usr/share/doc/cockpit-ws
/usr/share/doc/cockpit-ws/changelog.Debian.gz
/usr/share/doc/cockpit-ws/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/cockpit-ws
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cockpit-desktop.1.gz
/usr/share/man/man5
/usr/share/man/man5/cockpit.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/cockpit-tls.8.gz
/usr/share/man/man8/cockpit-ws.8.gz
/usr/share/man/man8/pam_ssh_add.8.gz
/usr/share/cockpit/branding/arch/apple-touch-icon.png
/usr/share/cockpit/branding/arch/favicon.ico
/usr/share/cockpit/branding/arch/logo.png
/usr/share/cockpit/branding/debian/favicon.ico
/usr/share/cockpit/branding/debian/logo.png
