# debianrfs/var/lib/dpkg/info/generic-sys-mods.list
/.
/etc
/etc/bbb.io
/etc/bbb.io/templates
/etc/bbb.io/templates/profile
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/61-smsc95xx-mac-address.rules
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bb_install_modules.service
/lib/systemd/system/grow_partition.service
/lib/systemd/system/regenerate_ssh_host_keys.service
/lib/systemd/system/resize_filesystem.service
/usr
/usr/bin
/usr/bin/bb-growpart
/usr/bin/bb-install-modules
/usr/bin/bb-resizefs
/usr/bin/bb-setup-mac-address
/usr/sbin
/usr/sbin/bb-smsc95xx-set-mac-address
/usr/share
/usr/share/doc
/usr/share/doc/generic-sys-mods
/usr/share/doc/generic-sys-mods/changelog.Debian.gz
