# debianrfs/var/lib/dpkg/info/gpsd.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.gpsd
/etc/default
/etc/default/gpsd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/gpsd.service
/lib/systemd/system/gpsd.socket
/lib/systemd/system/gpsdctl@.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-gpsd.rules
/usr
/usr/bin
/usr/bin/gpsctl
/usr/bin/ntpshmmon
/usr/bin/ppscheck
/usr/bin/ubxtool
/usr/sbin
/usr/sbin/gpsd
/usr/sbin/gpsdctl
/usr/share
/usr/share/doc
/usr/share/doc/gpsd
/usr/share/doc/gpsd/README.adoc.gz
/usr/share/doc/gpsd/TROUBLESHOOTING.adoc.gz
/usr/share/doc/gpsd/changelog.Debian.armhf.gz
/usr/share/doc/gpsd/changelog.Debian.gz
/usr/share/doc/gpsd/changelog.gz
/usr/share/doc/gpsd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/gpsctl.1.gz
/usr/share/man/man1/ntpshmmon.1.gz
/usr/share/man/man1/ubxtool.1.gz
/usr/share/man/man8
/usr/share/man/man8/gpsd.8.gz
/usr/share/man/man8/gpsdctl.8.gz
/usr/share/man/man8/gpsinit.8.gz
/usr/share/man/man8/ppscheck.8.gz
