# debianrfs/var/lib/dpkg/info/iwd.list
/.
/etc
/etc/init.d
/etc/init.d/iwd
/etc/iwd
/etc/iwd/main.conf
/lib
/lib/systemd
/lib/systemd/network
/lib/systemd/network/80-iwd.link
/lib/systemd/system
/lib/systemd/system/iwd.service
/usr
/usr/bin
/usr/bin/iwctl
/usr/bin/iwmon
/usr/lib
/usr/lib/modules-load.d
/usr/lib/modules-load.d/pkcs8.conf
/usr/libexec
/usr/libexec/iwd
/usr/share
/usr/share/bug
/usr/share/bug/iwd
/usr/share/bug/iwd/presubj
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/net.connman.iwd.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/iwd-dbus.conf
/usr/share/doc
/usr/share/doc/iwd
/usr/share/doc/iwd/README.Debian
/usr/share/doc/iwd/changelog.Debian.gz
/usr/share/doc/iwd/changelog.gz
/usr/share/doc/iwd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/iwctl.1.gz
/usr/share/man/man1/iwmon.1.gz
/usr/share/man/man5
/usr/share/man/man5/iwd.ap.5.gz
/usr/share/man/man5/iwd.config.5.gz
/usr/share/man/man5/iwd.network.5.gz
/usr/share/man/man7
/usr/share/man/man7/iwd.debug.7.gz
/usr/share/man/man8
/usr/share/man/man8/iwd.8.gz
