# debianrfs/var/lib/dpkg/info/mender-client.list
/.
/etc
/etc/mender
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/mender-client.service
/usr
/usr/bin
/usr/bin/mender
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/io.mender.AuthenticationManager.conf
/usr/share/dbus-1/system.d/io.mender.UpdateManager.conf
/usr/share/doc
/usr/share/doc/mender-client
/usr/share/doc/mender-client/README.Debian
/usr/share/doc/mender-client/changelog.Debian.armhf.gz
/usr/share/doc/mender-client/changelog.Debian.gz
/usr/share/doc/mender-client/copyright
/usr/share/doc/mender-client/examples
/usr/share/doc/mender-client/examples/directory-artifact-gen
/usr/share/doc/mender-client/examples/docker-artifact-gen
/usr/share/doc/mender-client/examples/mender-device-identity
/usr/share/doc/mender-client/examples/mender-inventory-bootloader-integration
/usr/share/doc/mender-client/examples/mender-inventory-geo
/usr/share/doc/mender-client/examples/mender-inventory-hostinfo
/usr/share/doc/mender-client/examples/mender-inventory-network
/usr/share/doc/mender-client/examples/mender-inventory-os
/usr/share/doc/mender-client/examples/mender-inventory-provides
/usr/share/doc/mender-client/examples/mender-inventory-rootfs-type
/usr/share/doc/mender-client/examples/mender-inventory-update-modules
/usr/share/doc/mender-client/examples/mender.conf.demo
/usr/share/doc/mender-client/examples/mender.conf.production
/usr/share/doc/mender-client/examples/single-file-artifact-gen
/usr/share/mender
/usr/share/mender/identity
/usr/share/mender/inventory
/usr/share/mender/modules
/usr/share/mender/modules/v3
/usr/share/mender/modules/v3/deb
/usr/share/mender/modules/v3/directory
/usr/share/mender/modules/v3/docker
/usr/share/mender/modules/v3/rpm
/usr/share/mender/modules/v3/script
/usr/share/mender/modules/v3/single-file
/var
/var/lib
/var/lib/mender
