# debianrfs/var/lib/dpkg/info/wireguard-tools.list
/.
/etc
/etc/wireguard
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/wg-quick.target
/lib/systemd/system/wg-quick@.service
/usr
/usr/bin
/usr/bin/wg
/usr/bin/wg-quick
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/wg
/usr/share/bash-completion/completions/wg-quick
/usr/share/doc
/usr/share/doc/wireguard-tools
/usr/share/doc/wireguard-tools/NEWS.Debian.gz
/usr/share/doc/wireguard-tools/README.Debian
/usr/share/doc/wireguard-tools/changelog.Debian.armhf.gz
/usr/share/doc/wireguard-tools/changelog.Debian.gz
/usr/share/doc/wireguard-tools/copyright
/usr/share/doc/wireguard-tools/examples
/usr/share/doc/wireguard-tools/examples/dns-hatchet
/usr/share/doc/wireguard-tools/examples/dns-hatchet/README
/usr/share/doc/wireguard-tools/examples/dns-hatchet/apply.sh
/usr/share/doc/wireguard-tools/examples/dns-hatchet/hatchet.bash
/usr/share/doc/wireguard-tools/examples/embeddable-wg-library
/usr/share/doc/wireguard-tools/examples/embeddable-wg-library/Makefile
/usr/share/doc/wireguard-tools/examples/embeddable-wg-library/README
/usr/share/doc/wireguard-tools/examples/embeddable-wg-library/test.c
/usr/share/doc/wireguard-tools/examples/embeddable-wg-library/wireguard.c
/usr/share/doc/wireguard-tools/examples/embeddable-wg-library/wireguard.h
/usr/share/doc/wireguard-tools/examples/extract-handshakes
/usr/share/doc/wireguard-tools/examples/extract-handshakes/Makefile
/usr/share/doc/wireguard-tools/examples/extract-handshakes/README
/usr/share/doc/wireguard-tools/examples/extract-handshakes/extract-handshakes.sh
/usr/share/doc/wireguard-tools/examples/extract-handshakes/offset-finder.c
/usr/share/doc/wireguard-tools/examples/extract-keys
/usr/share/doc/wireguard-tools/examples/extract-keys/Makefile
/usr/share/doc/wireguard-tools/examples/extract-keys/README
/usr/share/doc/wireguard-tools/examples/extract-keys/config.c
/usr/share/doc/wireguard-tools/examples/extract-keys/extract-keys.c
/usr/share/doc/wireguard-tools/examples/highlighter
/usr/share/doc/wireguard-tools/examples/highlighter/Makefile
/usr/share/doc/wireguard-tools/examples/highlighter/README
/usr/share/doc/wireguard-tools/examples/highlighter/fuzz.c
/usr/share/doc/wireguard-tools/examples/highlighter/gui
/usr/share/doc/wireguard-tools/examples/highlighter/gui/highlight.cpp
/usr/share/doc/wireguard-tools/examples/highlighter/gui/highlight.pro
/usr/share/doc/wireguard-tools/examples/highlighter/highlight.c
/usr/share/doc/wireguard-tools/examples/highlighter/highlighter.c
/usr/share/doc/wireguard-tools/examples/highlighter/highlighter.h
/usr/share/doc/wireguard-tools/examples/json
/usr/share/doc/wireguard-tools/examples/json/README
/usr/share/doc/wireguard-tools/examples/json/wg-json
/usr/share/doc/wireguard-tools/examples/keygen-html
/usr/share/doc/wireguard-tools/examples/keygen-html/README
/usr/share/doc/wireguard-tools/examples/keygen-html/keygen.html
/usr/share/doc/wireguard-tools/examples/keygen-html/wireguard.js
/usr/share/doc/wireguard-tools/examples/launchd
/usr/share/doc/wireguard-tools/examples/launchd/README
/usr/share/doc/wireguard-tools/examples/launchd/com.wireguard.wg0.plist
/usr/share/doc/wireguard-tools/examples/nat-hole-punching
/usr/share/doc/wireguard-tools/examples/nat-hole-punching/README
/usr/share/doc/wireguard-tools/examples/nat-hole-punching/nat-punch-client.c
/usr/share/doc/wireguard-tools/examples/nat-hole-punching/nat-punch-server.c
/usr/share/doc/wireguard-tools/examples/ncat-client-server
/usr/share/doc/wireguard-tools/examples/ncat-client-server/README
/usr/share/doc/wireguard-tools/examples/ncat-client-server/client-quick.sh
/usr/share/doc/wireguard-tools/examples/ncat-client-server/client.sh
/usr/share/doc/wireguard-tools/examples/ncat-client-server/server.sh
/usr/share/doc/wireguard-tools/examples/reresolve-dns
/usr/share/doc/wireguard-tools/examples/reresolve-dns/README
/usr/share/doc/wireguard-tools/examples/reresolve-dns/reresolve-dns.sh
/usr/share/doc/wireguard-tools/examples/sticky-sockets
/usr/share/doc/wireguard-tools/examples/sticky-sockets/README
/usr/share/doc/wireguard-tools/examples/sticky-sockets/sticky-sockets.c
/usr/share/doc/wireguard-tools/examples/synergy
/usr/share/doc/wireguard-tools/examples/synergy/README
/usr/share/doc/wireguard-tools/examples/synergy/synergy-client.sh
/usr/share/doc/wireguard-tools/examples/synergy/synergy-server.sh
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/wireguard-tools
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/wg-quick.8.gz
/usr/share/man/man8/wg.8.gz
