# rfs/var/lib/dpkg/info/live-config-upstart.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/live-config-upstart
/usr/share/doc/live-config-upstart/copyright
/usr/share/doc/live-config-upstart/changelog.Debian.gz
/etc
/etc/init.d
/etc/init
/etc/init/live-config.conf
/etc/init.d/live-config
